# package/vpntfi/buildroot/etc/systemd/system/radiusd.service
[Unit]
Description=freeradius service

[Service]
Type=forking
ExecStart=/usr/sbin/radiusd -t

[Install]
WantedBy=multi-user.target
